# .config/systemd/user/restic.timer
[Unit]
Description=Scheduled daily backup

[Timer]
OnCalendar=Mon-Fri *-*-* 12:00
Unit=restic.service

[Install]
WantedBy=timers.target
